# rootfs/etc/services.d/ha-nostr-alert/finish
#!/command/with-contenv bashio
# ==============================================================================
# HA Nostr Alert Service Finisher
# ==============================================================================

echo "=== HA Nostr Alert Service Finished ==="
